# systemd/meshmini.service
[Unit]
Description=MeshMini - minimal Meshtastic BBS (with peer sync)
After=network.target

[Service]
Type=simple
User=pihole
Group=pihole
WorkingDirectory=/opt/meshmini
Environment=PYTHONUNBUFFERED=1
Environment=MMB_DB=/opt/meshmini/board.db
Environment=MMB_DEVICE=auto
Environment=MMB_NAME=MeshLink BBS
Environment=MMB_ADMINS=!deadbeef
Environment=MMB_PEERS=!cafef00d,!baadf00d
ExecStart=/opt/meshmini/venv/bin/python /opt/meshmini/meshmini.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
